# conf/crontab
1 * * * * scripts/store_counts.py conf/infobase.yaml conf/openlibrary.yaml conf/coverstore.yaml 1
