# conf/crontab
1 * * * * scripts/store_counts.py conf/infobase.yml conf/openlibrary.yml 1
