# conf/supervisor/macosx.ini
# Services run only on mac os x

[program:postgres]
command=vendor/postgresql-8.4.4/bin/postgres -D var/lib/postgresql
priority=100
stdout_logfile = var/log/postrgesql.log
redirect_stderr = true
